# detection-rules/sigma/persistence/registry_run_keys.yml
title: Suspicious Registry Run Key Modification
id: 34567890-3456-3456-3456-345678901cde
status: stable
description: Detects modifications to registry run keys often used for persistence mechanisms
references:
    - https://attack.mitre.org/techniques/T1547/001/
    - https://docs.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys
author: Defensive Toolkit
date: 2025-10-15
modified: 2025-10-15
tags:
    - attack.persistence
    - attack.t1547.001
    - attack.privilege_escalation
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        TargetObject|contains:
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Policies\Explorer\Run\'
    filter_legitimate:
        Image|startswith:
            - 'C:\Program Files\'
            - 'C:\Program Files (x86)\'
            - 'C:\Windows\System32\'
    condition: selection and not filter_legitimate
falsepositives:
    - Legitimate software installations
    - System updates
    - Enterprise software deployment
level: medium
